# debian/prerm
#!/bin/bash
systemctl stop wildfly || true
systemctl disable wildfly || true

