# debian/prerm
#!/bin/bash
systemctl stop wildfly
systemctl disable wildfly

